# config/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
logfile = /var/log/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = False
minfds = 1024
minprocs = 200
umask = 022
identifier = supervisor
directory = /tmp
nocleanup = true
childlogdir = /tmp

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:pinggo]
command=%(ENV_HOME)s/goworkspace/src/go_ping_handler/ping
autostart=false
redirect_stderr=true
autorestart=true
startretries=10
directory=%(ENV_HOME)s/goworkspace/src/go_ping_handler
user=deployer
stdout_logfile=%(ENV_HOME)s/goworkspace/log/ping_supervisor.log
stderr_logfile=%(ENV_HOME)s/goworkspace/log/ping_supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=GOPATH="%(ENV_HOME)s/goworkspace/",GO_ENV="var_go_env"

[program:pinglistener]
command=%(ENV_HOME)s/goworkspace/src/go_ping_handler/ping_listener
autostart=false
autorestart=true
redirect_stderr=true
autorestart=true
startretries=10
directory=%(ENV_HOME)s/goworkspace/src/go_ping_handler
user=deployer
stdout_logfile=%(ENV_HOME)s/goworkspace/log/listener_supervisor.log
stderr_logfile=%(ENV_HOME)s/goworkspace/log/listener_supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=GOPATH="%(ENV_HOME)s/goworkspace/",GO_ENV="var_go_env"

[program:notification]
command=%(ENV_HOME)s/goworkspace/src/go_notification/notification
autostart=false
redirect_stderr=true
autorestart=true
startretries=10
directory=%(ENV_HOME)s/goworkspace/src/go_notification
user=deployer
stdout_logfile=%(ENV_HOME)s/goworkspace/log/notification_supervisor.log
stderr_logfile=%(ENV_HOME)s/goworkspace/log/notification_supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=GOPATH="%(ENV_HOME)s/goworkspace/",GO_ENV="var_go_env"

[program:s2_wrapper]
command=%(ENV_HOME)s/goworkspace/src/s2_wrapper/main
autostart=false
redirect_stderr=true
autorestart=true
startretries=10
directory=%(ENV_HOME)s/goworkspace/src/s2_wrapper
user=deployer
stdout_logfile=%(ENV_HOME)s/goworkspace/log/s2_supervisor.log
stderr_logfile=%(ENV_HOME)s/goworkspace/log/s2_supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=GOPATH="%(ENV_HOME)s/goworkspace/",GO_ENV="var_go_env"


[program:manager]
command=%(ENV_HOME)s/goworkspace/src/serviceability_manager/manager
autostart=false
redirect_stderr=true
autorestart=true
startretries=10
directory=%(ENV_HOME)s/goworkspace/src/serviceability_manager
stdout_logfile=%(ENV_HOME)s/goworkspace/log/serviceability_manager_supervisor.log
stderr_logfile=%(ENV_HOME)s/goworkspace/log/serviceability_manager_supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=GOPATH="%(ENV_HOME)s/goworkspace/",GO_ENV="var_go_env"

[program:serviceability_consumer]
command=%(ENV_HOME)s/goworkspace/src/serviceability_manager/serviceability_consumer
autostart=false
redirect_stderr=true
autorestart=true
startretries=10
directory=%(ENV_HOME)s/goworkspace/src/serviceability_manager
stdout_logfile=%(ENV_HOME)s/goworkspace/log/serviceability_consumer_supervisor.log
stderr_logfile=%(ENV_HOME)s/goworkspace/log/serviceability_consumer_supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=GOPATH="%(ENV_HOME)s/goworkspace/",GO_ENV="var_go_env"

[program:serviceability_producer]
command=%(ENV_HOME)s/goworkspace/src/serviceability_manager/serviceability_producer
autostart=false
redirect_stderr=true
autorestart=true
startretries=10
directory=%(ENV_HOME)s/goworkspace/src/serviceability_manager
stdout_logfile=%(ENV_HOME)s/goworkspace/log/serviceability_producer_supervisor.log
stderr_logfile=%(ENV_HOME)s/goworkspace/log/serviceability_producer_supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=GOPATH="%(ENV_HOME)s/goworkspace/",GO_ENV="var_go_env"

[program:serviceability_pusher]
command=%(ENV_HOME)s/goworkspace/src/serviceability_manager/serviceability_pusher
autostart=false
redirect_stderr=true
autorestart=true
startretries=10
directory=%(ENV_HOME)s/goworkspace/src/serviceability_manager
stdout_logfile=%(ENV_HOME)s/goworkspace/log/serviceability_pusher_supervisor.log
stderr_logfile=%(ENV_HOME)s/goworkspace/log/serviceability_pusher_supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=GOPATH="%(ENV_HOME)s/goworkspace/",GO_ENV="var_go_env"

[program:serviceability_monitor]
command=%(ENV_HOME)s/goworkspace/src/serviceability_manager/serviceability_monitor
autostart=false
redirect_stderr=true
autorestart=true
startretries=10
directory=%(ENV_HOME)s/goworkspace/src/serviceability_manager
stdout_logfile=%(ENV_HOME)s/goworkspace/log/serviceability_monitor_supervisor.log
stderr_logfile=%(ENV_HOME)s/goworkspace/log/serviceability_monitor_supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=GOPATH="%(ENV_HOME)s/goworkspace/",GO_ENV="var_go_env"
